# curator2-1243875156.conf
[program:curator2-1243875156]
command=python3 /home/curator2/curator2-1243875156/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator2/curator2-1243875156/test.err.log
stdout_logfile=/home/curator2/curator2-1243875156/test.out.log